# board/rasputin/common/rootfs-overlay/etc/systemd/system/systemd-networkd-wait-online.service.d/any.conf
# Wait for ANY link to be routable, not ALL. An appliance may ship with several
# NICs (the n100 has i225/e1000e/r8169) but only one cabled; the default
# wait-online waits for EVERY matched link, so an unplugged NIC blocks boot for
# the full ~2min timeout — stalling first boot and delaying the console IP banner
# (rasputin-issue-ip runs After=network-online.target). --any proceeds as soon as
# one link is online. See geekdojo/rasputin-os#10.
[Service]
ExecStart=
ExecStart=/usr/lib/systemd/systemd-networkd-wait-online --any
