# board/rasputin/common/rootfs-overlay/etc/systemd/system/serial-getty@ttyAMA0.service.d/10-rasputin-bmc-host.conf
# On a BMC-host node, serial0/ttyAMA0 is the command channel to the BMC
# bus — a login getty on it would interleave with BMC replies and console
# bridges (control-plane/bmc-bitscope.md §5). The marker is written by
# firstboot from the seed's RASPUTIN_BMC_HOST=1; /etc is read-only
# squashfs, so this condition (not runtime masking) is the gate — the
# same idiom as role.controlplane. Covers both the Buildroot-enabled
# instance and the one systemd-getty-generator spawns from console=.
[Unit]
ConditionPathExists=!/var/lib/rasputin/bmc-host
