# board/rasputin/common/rootfs-overlay/etc/systemd/system/serial-getty@ttyS0.service.d/10-rasputin-bmc-host.conf
# Same gate as the ttyAMA0 drop-in: no serial login getty on a BMC-host
# node (control-plane/bmc-bitscope.md §5). The BitScope host is Pi-only
# today, but the policy is generic — a future chassis BMC host may be
# amd64, where the serial console is ttyS0.
[Unit]
ConditionPathExists=!/var/lib/rasputin/bmc-host
