# board/rasputin/common/rootfs-overlay/etc/systemd/system/dropbear.service
# Key-only SSH for support/debugging a headless controlplane. rasputin-os has
# no display attached in normal operation, and until the mesh forms there's no
# tailnet path in — so dropbear is the supported remote shell.
#
# Posture: -s disables password auth entirely, so even though console root is
# passwordless (bench), the NETWORK vector is key-only. With no authorized key
# present, nobody can SSH in (console still works) — never passwordless root
# over the network.
#
# Authorized keys live on the PERSISTENT partition (the rootfs is read-only
# squashfs, so /root/.ssh can't hold runtime-supplied keys): dropbear runs
# with -D /var/lib/rasputin/dropbear and reads authorized_keys there. The
# ONLY source is the seed — RASPUTIN_SSH_AUTHORIZED_KEY in rasputin-seed.env,
# merged append-if-missing by rasputin-firstboot (merge-authorized-keys.sh).
# The image bakes NO key and the build cannot inject one (pre-GA vendor-key
# removal, 2026-07-09) — the same path serves us and end users. The file is
# operator-owned: manual additions survive, revocation is a manual edit.
# Host keys live alongside for the same read-only-rootfs reason (stable host
# identity across reboots and A/B updates); they're generated on first boot
# if absent.
[Unit]
Description=Dropbear SSH server (key-only)
After=network.target systemd-tmpfiles-setup.service rasputin-firstboot.service
Wants=network.target
RequiresMountsFor=/var/lib/rasputin

[Service]
Type=simple
ExecStartPre=/bin/sh -c 'mkdir -p /var/lib/rasputin/dropbear && { test -s /var/lib/rasputin/dropbear/ed25519_host_key || /usr/bin/dropbearkey -t ed25519 -f /var/lib/rasputin/dropbear/ed25519_host_key; }'
# dropbear re-reads authorized_keys on every auth attempt, so a key merged
# after startup (firstboot's seed key, a manual append) needs no restart.
# -F foreground (systemd-managed), -E log to the journal, -s NO password auth
# (key-only). -r pins the persistent host key; -D makes dropbear read
# <dir>/authorized_keys instead of ~/.ssh/authorized_keys.
ExecStart=/usr/sbin/dropbear -F -E -s -r /var/lib/rasputin/dropbear/ed25519_host_key -D /var/lib/rasputin/dropbear
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
