# board/rasputin/common/rootfs-overlay/etc/systemd/system/rasputin-api.service
[Unit]
Description=Rasputin control-plane API
After=network-online.target rasputin-firstboot.service
Wants=network-online.target
Requires=rasputin-firstboot.service
RequiresMountsFor=/var/lib/rasputin
# Role gate: this unit ships preset-enabled on every image but only starts
# on the controlplane — rasputin-firstboot drops this marker on the
# persistent partition when role==controlplane. (Not ConditionEnvironment:
# that checks PID 1's environment, which never contains node.env. Not
# runtime `systemctl enable`: /etc is read-only squashfs.)
# See os-images/provisioning.md §2.
ConditionPathExists=/var/lib/rasputin/role.controlplane

[Service]
# Image defaults via Environment=; the seed-written node.env can OVERRIDE
# any of them — systemd.exec(5): "Settings from these files [EnvironmentFile]
# override settings made with Environment=", and they're ordered after the
# Environment= lines here as belt-and-braces anyway.
Environment=RASPUTIN_DATA_DIR=/var/lib/rasputin
# First-run surface: HTTP :80 is the bootstrap/redirect listener (and keeps
# serving /healthz over plain HTTP — the QEMU smoke and the agent probe it);
# HTTPS :443 carries the real UI/API. The unit runs as root, so low ports
# need no extra capability grants. WebAuthn RP ID is the mDNS name the
# image announces (rasputin-hostname.service + resolved): rasputin.local.
Environment=RASPUTIN_HTTP_ADDR=:80
Environment=RASPUTIN_HTTPS_ADDR=:443
Environment=RASPUTIN_RP_ID=rasputin.local
Environment=RASPUTIN_RP_ORIGINS=https://rasputin.local
Environment=RASPUTIN_PUBLIC_BASE_URL=https://rasputin.local
# Embedded NATS must be reachable from the other nodes' agents (firewall,
# compute) — they attach LAN-first, mesh-second (first-run-bootstrap.md
# decision #6). The api's code default is 127.0.0.1 (dev posture).
# ⚠ MVS trust posture: the bus is currently TOFU/unauthenticated
# (architecture.md decision #14 — production NKey/JWT via the seed join
# token is the gating item before design partners). LAN exposure here is
# the documented interim; override to 127.0.0.1 via node.env to close it.
Environment=RASPUTIN_NATS_HOST=0.0.0.0
# Self-hosted Headscale advertises THIS as its server_url — the address every
# node's tailscaled connects to. Pin the stable mDNS name, not the controlplane's
# DHCP IP: resolveServerHost() would otherwise bake the current lease into
# Headscale + every node's tailscale state, so a lease change drops the whole
# mesh (bench: .180→.181→.182 across reboots). The leaf cert already carries
# rasputin.local in its SAN. URL-only (no RASPUTIN_HEADSCALE_API_KEY) keeps the
# backend self-hosted — external mode needs BOTH (main.go wireMesh hasExternal).
# Requires every node to resolve rasputin.local at the OS level (CP: resolved
# mDNS via the stub resolv.conf below; firewall: its own mDNS resolver).
Environment=RASPUTIN_HEADSCALE_URL=https://rasputin.local:18080
EnvironmentFile=-/var/lib/rasputin/node.env
ExecStart=/usr/bin/rasputin-api
Restart=always
RestartSec=2
# Watchdog contract: the api pets via sd_notify(WATCHDOG=1) gated on a
# SQLite liveness probe — implemented in control-plane v0.4.3
# (api/internal/sdnotify). BEFORE that the binary never notified and this
# stanza silently SIGABRT-killed the api every ~33s on every hardware
# boot (found 2026-06-12; the QEMU smoke's retry-until-success healthz
# polls fit inside the live windows, so five releases shipped with it).
# Type=notify makes the contract explicit — and the smoke now asserts
# 70s of uninterrupted uptime so this class of bug can't ship again.
Type=notify
WatchdogSec=30
NotifyAccess=main

[Install]
WantedBy=multi-user.target
