# board/rasputin/common/rootfs-overlay/etc/systemd/system/rasputin-agent.service
[Unit]
Description=Rasputin node agent
After=network-online.target rasputin-firstboot.service
Wants=network-online.target
Requires=rasputin-firstboot.service
RequiresMountsFor=/var/lib/rasputin

[Service]
# Agent state (app compose files, updater/mesh backend state) must land on
# the persistent partition — the rootfs is read-only squashfs. The agent
# uses this value verbatim (no per-node suffix; one agent per host).
# node.env can still override it: EnvironmentFile beats Environment=.
Environment=RASPUTIN_AGENT_STATE_DIR=/var/lib/rasputin/agent-state
# node.env is written by rasputin-firstboot onto the persistent partition.
EnvironmentFile=-/var/lib/rasputin/node.env
# Safety net for anything else that resolves paths relative to cwd (the
# env var above is the real fix for the agent's state dir).
WorkingDirectory=/var/lib/rasputin
ExecStart=/usr/bin/rasputin-agent
Restart=always
RestartSec=2
# systemd watchdog — part of the update defense-in-depth (the CM5
# bootloader can't count boot attempts, so userspace must catch a wedge).
# The agent pets via sd_notify(WATCHDOG=1) — implemented in control-plane
# v0.4.3 (agent/internal/sdnotify); BEFORE that the binary never notified
# and this stanza silently SIGABRT-killed the agent every 30s on every
# hardware boot since the scaffold (found 2026-06-12 on the Mu bench —
# every probe we had fit inside the ~28s live windows). Type=notify makes
# the contract explicit: a binary without notify support now fails START
# loudly instead of looping quietly. Do not pair WatchdogSec with a
# non-notifying binary again. See os-images/buildroot-os.md §3.
Type=notify
WatchdogSec=30
NotifyAccess=main

[Install]
WantedBy=multi-user.target
