# board/rasputin/common/rootfs-overlay/etc/systemd/system/rasputin-mark-good.service
# Marks the currently-booted RAUC slot "good" once the OS has reached a working
# userspace, to reset the GRUB boot-counter. n100 (GRUB) ONLY.
#
# Why n100-only: on the n100 the grub.cfg boot-counter consumes one try per
# boot; without a mark-good on each healthy boot, a normal reboot would see
# TRY=1 and fall back to the other slot (defense-in-depth LAYER 1: the OS
# booted). RAUC's grubenv carries a persistent OK flag, so the update saga can
# still `mark-bad` AFTER this to force a rollback (mark-bad → OK=0, which grub
# re-evaluates every boot — this auto-good doesn't prevent rollback).
#
# The Pi (custom tryboot backend) is the OPPOSITE: there marking the running
# slot good is the COMMIT (promotes autoboot.txt's [all] boot_partition) — an
# irreversible point with no boot-counter to fall back on. Auto-committing on
# OS-up would pre-empt the update saga's health check, so a later mark-bad
# couldn't roll back. So on the Pi this service is gated OFF, and the saga's
# health-gated mark-good (api healthCheckAndCommit → UpdateMarkGoodCmd) is the
# sole committer — giving the Pi full app-health rollback within the trial
# window. The gate is the presence of the Pi selector's autoboot.txt (the same
# marker the agent's rebootArgs + rasputin-rauc-reconcile.service key off).
#
# Ordered after the agent so "good" means we reached our own userspace; after
# the seed mount so the autoboot.txt gate is evaluated against the mounted FAT;
# not coupled to network so a network blip can't trigger a spurious rollback.
# Skips silently (no failure) on a node without RAUC config, or on any Pi.

[Unit]
Description=Mark the running RAUC slot good — reset GRUB boot-counter (n100 only)
After=rasputin-agent.service run-rasputin\x2dseed.mount
ConditionPathExists=/etc/rauc/system.conf
# n100-only: skip on the Pi, whose selector FAT carries autoboot.txt (the Pi
# commit is saga-driven, not auto-on-OS-up — see the header).
ConditionPathExists=!/run/rasputin-seed/autoboot.txt

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/rauc status mark-good

[Install]
WantedBy=multi-user.target
