# board/rasputin/common/rootfs-overlay/etc/systemd/system/rasputin-issue-ip.service
[Unit]
Description=Refresh the console login banner with the node's IP
# /etc/issue carries `\4` (IPv4); agetty resolves it when it renders the issue,
# which happens at getty start — often BEFORE DHCP has assigned an address, so
# the IP line would otherwise show blank. Once the network is up, `agetty
# --reload` (it touches /run/agetty.reload) makes every running agetty
# re-display the issue with the now-assigned address. Read-only /etc means we
# can't rewrite /etc/issue at runtime, so `\4` + reload is the read-only-
# friendly way to surface a live IP on the console login screen.
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/agetty --reload

[Install]
WantedBy=multi-user.target
