# board/rasputin/common/rootfs-overlay/etc/systemd/system/rasputin-agent.service
[Unit]
Description=Rasputin node agent
After=network-online.target rasputin-firstboot.service
Wants=network-online.target
Requires=rasputin-firstboot.service
RequiresMountsFor=/var/lib/rasputin

[Service]
# Agent state (app compose files, updater/mesh backend state) must land on
# the persistent partition — the rootfs is read-only squashfs. The agent
# uses this value verbatim (no per-node suffix; one agent per host).
# node.env can still override it: EnvironmentFile beats Environment=.
Environment=RASPUTIN_AGENT_STATE_DIR=/var/lib/rasputin/agent-state
# node.env is written by rasputin-firstboot onto the persistent partition.
EnvironmentFile=-/var/lib/rasputin/node.env
# Safety net for anything else that resolves paths relative to cwd (the
# env var above is the real fix for the agent's state dir).
WorkingDirectory=/var/lib/rasputin
ExecStart=/usr/bin/rasputin-agent
Restart=always
RestartSec=2
# systemd watchdog — part of the update defense-in-depth (the CM5
# bootloader can't count boot attempts, so userspace must catch a wedge).
# The agent calls sd_notify(WATCHDOG=1); if it stops pinging, systemd
# kills + restarts it, and after enough failures the bootloader retries
# the other slot. See os-images/buildroot-os.md §3.
WatchdogSec=30
NotifyAccess=main

[Install]
WantedBy=multi-user.target
