# board/rasputin/common/rootfs-overlay/etc/systemd/system/rasputin-firstboot.service
[Unit]
Description=Rasputin first-boot provisioning (read seed, set role/identity)
# Wants (not Requires) the seed mount: a missing/unmountable seed degrades to
# the script's documented fallbacks (kernel cmdline, then defaults) instead of
# blocking provisioning outright. provisioning.md §1.
Wants=run-rasputin\x2dseed.mount
After=local-fs.target run-rasputin\x2dseed.mount
# Everything this script writes lands on the persistent partition; if that
# mount failed (nofail in fstab), provisioning must not half-run.
RequiresMountsFor=/var/lib/rasputin
Before=rasputin-agent.service rasputin-api.service
ConditionPathExists=!/var/lib/rasputin/.provisioned

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/rasputin/firstboot/rasputin-firstboot.sh

[Install]
WantedBy=multi-user.target
