# board/rasputin/common/rootfs-overlay/etc/systemd/system/rasputin-issue-ip.service
[Unit]
Description=Write the node's IP into the console login banner
# /etc/issue is a symlink to /run/issue (writable tmpfs); this oneshot fills it
# with the node's REAL routable IPv4 (default-route source addr) and reloads
# agetty so the login screen shows it. Reads the route source rather than
# agetty's `\4` so it never flashes loopback/link-local before DHCP settles.
# See usr/lib/rasputin/issue/rasputin-issue-ip.sh.
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/rasputin/issue/rasputin-issue-ip.sh

[Install]
WantedBy=multi-user.target
