# board/rasputin/common/rootfs-overlay/etc/systemd/system/run-rasputin\x2dseed.mount
# Mounts the boot FAT (the provisioning seed lives there) at /run/rasputin-seed
# so rasputin-firstboot.sh can read rasputin-seed.env. os-images/provisioning.md §1.
#
# Mounted by FILESYSTEM LABEL, not PARTLABEL: the GPT partition is named "esp"
# on n100 and "firmware" on cm5 (see each board's genimage.cfg), but both FATs
# carry the same volume label RASPUTIN-FW — one unit serves both boards.
#
# Read-write: RAUC's grub backend updates grubenv (/EFI/BOOT/grubenv on this
# FAT) at runtime on mark-good/mark-bad/activate, so the boot FAT must be
# writable. The only runtime writer is grub-editenv (a fixed ~1 KB in-place
# block); firstboot still only reads the seed. (Was ro; widened for the n100
# RAUC GRUB boot-counter — os-images/buildroot-os.md §3.)
#
# Filename is the systemd-escape of /run/rasputin-seed (the "\x2d" IS the
# escaped dash, literal in the filename — do not "fix" it).

[Unit]
Description=Rasputin provisioning seed + grubenv (boot FAT)
DefaultDependencies=no
Before=umount.target
Conflicts=umount.target

[Mount]
What=/dev/disk/by-label/RASPUTIN-FW
Where=/run/rasputin-seed
Type=vfat
Options=rw
DirectoryMode=0755
