# board/rasputin/common/rootfs-overlay/etc/systemd/system/tailscaled.service.d/10-rasputin-mesh-ca.conf
# Mesh CA trust for tailscaled. Rasputin runs a self-hosted Headscale control
# server with a per-installation Mesh CA (not a public root). tailscaled is a
# Go daemon; crypto/x509 reads SSL_CERT_FILE in ADDITION to the system cert
# dirs, so this points it at the Mesh CA bundle the agent writes on enrollment
# while leaving the public roots (DERP relays) trusted. Path is the agent's
# RASPUTIN_MESH_CA_BUNDLE default, so no agent env change is needed.
#
# A drop-in (not an edit of the packaged unit) so a package rebuild can't
# clobber it.
[Service]
Environment=SSL_CERT_FILE=/var/lib/rasputin/mesh/tailscaled-ca.pem
