# board/rasputin/common/rootfs-overlay/etc/systemd/system/tailscaled.service.d/persistent-state.conf
# /var/lib/tailscale is a baked-in symlink to /var/lib/rasputin/tailscale (the
# rootfs is read-only squashfs; tailscaled creates its state dir by following
# the symlink onto the persistent partition). Same pattern as docker/containerd.
# Don't race the persistent mount.
[Unit]
RequiresMountsFor=/var/lib/rasputin

# Clear the upstream unit's CacheDirectory=tailscale. systemd creates a
# CacheDirectory under /var/cache, which is on the read-only squashfs — so the
# directive fails at step CACHE_DIRECTORY with EROFS and tailscaled never even
# execs (status=239/CACHE_DIRECTORY; the service then trips its start-limit and
# stays failed — bench 2026-06-18, dev.20: the controlplane's own tailscaled
# was down so it couldn't self-enroll). StateDirectory works only because
# /var/lib/tailscale is the symlink above; /var/cache has no such redirect.
# tailscaled doesn't need a persistent cache (it re-fetches the DERP map on
# start), so clear the directive rather than redirect it. An empty assignment
# resets the inherited value.
[Service]
CacheDirectory=
