# board/rasputin/common/rootfs-overlay/etc/systemd/system/tailscaled.service.d/persistent-state.conf
# /var/lib/tailscale is a baked-in symlink to /var/lib/rasputin/tailscale (the
# rootfs is read-only squashfs; tailscaled creates its state dir by following
# the symlink onto the persistent partition). Same pattern as docker/containerd.
# Don't race the persistent mount.
[Unit]
RequiresMountsFor=/var/lib/rasputin
