# board/rasputin/common/rootfs-overlay/etc/systemd/system/containerd.service.d/persistent-state.conf
# /var/lib/containerd is a baked-in symlink to /var/lib/rasputin/containerd
# (read-only squashfs rootfs). Don't race the persistent mount. Same pattern
# as docker.service.d/persistent-state.conf — docker found dead on dev.4
# because containerd silently had the identical problem one layer down.
[Unit]
RequiresMountsFor=/var/lib/rasputin
