# board/rasputin/common/rootfs-overlay/etc/systemd/system/docker.service.d/persistent-state.conf
# /var/lib/docker is a baked-in symlink to /var/lib/rasputin/docker (the
# rootfs is read-only squashfs; dockerd's MkdirAll follows the symlink and
# creates the target on the persistent partition). Don't let docker race
# the mount.
[Unit]
RequiresMountsFor=/var/lib/rasputin
